# updates/1.1.2/home/ethos/.config/autostart/ethos-overclock.desktop
[Desktop Entry]
Type=Application
Exec=sudo /opt/ethos/sbin/ethos-overclock boot
