# updates/1.0.2/home/ethos/.config/autostart/ethos-terminal.desktop
[Desktop Entry]
Type=Application
Exec=/opt/ethos/sbin/ethos-terminal 
