# updates/1.0.2/etc/init/ethos-watchdog.conf
# Ethos upstart job at /etc/init/ethos-watchdog.conf

pre-start script
               echo "ethOS-watchdog Started.."
end script

start on starting lightdm
stop on runlevel shutdown
respawn
respawn limit 10 90

script
 /opt/ethos/sbin/ethos-watchdog
end script

