# updates/1.0.4/etc/init/ethos-stats-daemon.conf
# Ethos upstart job at /etc/init/ethos-stats-daemon.conf

pre-start script
               echo "ethOS-stats-daemon Started.."
end script

start on starting lightdm
stop on runlevel shutdown
respawn
respawn limit 15 90

script
 /usr/bin/php /opt/ethos/sbin/ethos-stats-daemon
end script

post-start script
		/opt/ethos/sbin/ethos-motd-generator
		/usr/sbin/dmidecode | grep Product | tail -1 | cut -d":" -f2 | xargs > /tmp/motherboard.file
		/usr/bin/aticonfig --lsa | egrep 'AMD|67B1|6939|7300|Supported device' | wc -l > /tmp/gpucount.file
		/usr/bin/lspci | grep VGA | cut -d"." -f1,3 | sed 's#\[AMD/ATI\]##g' | sed 's/:00. //g' | sed 's/Radeon //g' | sed 's/OEM //g' > /tmp/gpulistconky.file
		/usr/bin/lspci | grep VGA | cut -d":" -f1,3 > /tmp/gpulist.file
		
end script

