# updates/1.1.2/etc/init/ethos-miner-monitor.conf
# Ethos upstart job at /etc/init/ethos-miner-monitor.conf

pre-start script
               echo "ethOS-miner-monitor Started.."
end script

start on starting lightdm
stop on runlevel shutdown
respawn
respawn limit 10 90

script
 /opt/ethos/sbin/ethos-miner-monitor
end script

